# share/docker/etc/services.d/11-sshd/run
#!/bin/sh -e

echo "[services.d] starting service: sshd"

/usr/sbin/sshd -4D
